# mayhem-for-api/supervisord.conf
[supervisord]
nodaemon=true

[program:quake-server]
command=/opt/quake/cmd/quake-server/quake-server

[program:rest-proxy]
command=/opt/quake/mayhem-for-api/rest-proxy